# mandrel-health.timer
[Unit]
Description=Run Mandrel Health Check every 30 seconds
Requires=mandrel.service

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
