# deploy/systemd/mandrel-command.service.d/onfailure-alert.conf
# Drop-in: fire ridge-failure-alert@mandrel-command.service when
# mandrel-command.service fails. %n expands to the failing unit name, so the
# alert template receives the exact unit that failed. Additive — overrides
# nothing else.
[Unit]
OnFailure=ridge-failure-alert@%n.service
