# deploy/systemd/mandrel.service.d/onfailure-alert.conf
# Drop-in: fire ridge-failure-alert@mandrel.service when mandrel.service fails.
# %n expands to the failing unit name (mandrel.service), so the alert template
# receives the exact unit that failed. Additive — overrides nothing else.
[Unit]
OnFailure=ridge-failure-alert@%n.service
